# crontab
# Collect temperature data every 5 minutes
*/5 * * * * root cd /app && python collect_data.py >> /var/log/cron.log 2>&1

# Empty line at end required for crontab